# goTrack.service
[Unit]
Description=goTrack service to activate at boot

[Service]
ExecStart=/usr/local/bin/goTrack

[Install]
WantedBy=multi-user.target